# crons/rtc-i2c.service
[Unit]
ConditionCapability=CAP_SYS_TIME
ConditionVirtualization=!container
DefaultDependencies=no
Wants=dev-rtc_i2c.device
After=dev-rtc_i2c.device
Before=system-timesyncd.service ntpd.service chrony.service

[Service]
Type=oneshot
CapabilityBoundingSet=CAP_SYS_TIME
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
DeviceAllow=/dev/rtc_i2c rw
DevicePolicy=closed
ExecStart=/sbin/hwclock --rtc=/dev/rtc_i2c --hctosys
ExecStartPost=/sbin/logger --tag rtc-i2c 'System clock set from i2c hardware clock'

[Install]
WantedBy=time-sync.target